# migration_systemd_services/migration_retry.service
[Unit]
Description=Migration Retry Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/migration/bin/pythonMigration /opt/migration/bin/migration_retry.py
StandardOutput=journal+console
StandardError=inherit

[Install]
WantedBy=multi-user.target
