# migration_systemd_services/migration_manager.service
[Unit]
Description=Migration Manager Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/migration/bin/pythonMigration /opt/migration/bin/migration_manager.py
StandardOutput=journal+console
StandardError=inherit

[Install]
WantedBy=multi-user.target
